# shell/.config/systemd/user/terminal.service
[Unit]
Description=Terminal
PartOf=graphical-session.target

[Service]
ExecStart=/usr/bin/xfce4-terminal

[Install]
WantedBy=xsession.target
